# debian/cron.weekly
#! /bin/sh
#
# cron script to rotate the vbox logfiles.

cd /var/log/isdn
for LOG in vboxgetty-ttyI[0-9].log 
do
        if [ -f $LOG ]; then
		savelog $LOG >/dev/null
	fi
done

killall vboxgetty 

exit 0
