# debian/isdnvboxserver.conffiles
/etc/isdn/vboxd.conf
/etc/isdn/vboxgetty.conf
/etc/cron.weekly/isdnvboxserver
