# debian/isdnutils-base.conffiles
/etc/init.d/isdnutils
/etc/isdn/init.d.functions
