# debian/isdneurofile.conffiles
/etc/init.d/isdneurofile
/etc/isdn/eft.conf
